# infra/deploy/systemd/searchbarbara-stable.template.service
[Unit]
# Repo template only.
# The installed stable unit may live under /etc/systemd/system/searchbarbara-stable.service.
# Update paths, user, and env file before installing this template.
Description=SearchBarbara Web Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/SearchBarbara
EnvironmentFile=/root/SearchBarbara/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/root/SearchBarbara/.venv/bin/python /root/SearchBarbara/run_web.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
